# recipes-core/server/lcdstatus-0.1/lcdstatus.service
[Unit]
Description=LCD Status

[Service] 
Type=simple
ExecStart=/sbin/lcdstatus

[Install]
WantedBy=basic.target
