# recipes-core/packages/ntpserver-abctl/ntpserver-boot-ok.service
[Unit]
Description=Mark new A/B slot as healthy after successful boot
After=multi-user.target
Wants=multi-user.target
ConditionPathExists=/boot/grub/grubenv

[Service]
Type=oneshot
ExecStart=/usr/sbin/ntpserver-boot-ok.sh

[Install]
WantedBy=multi-user.target
