# recipes-core/packages/lcdstatus-0.1/lcdstatus.service
[Unit]
Description=LCD Status

[Service] 
Type=simple
ExecStart=/usr/bin/lcdstatus

[Install]
WantedBy=basic.target
