# recipes-core/packages/abctl/boot-ok.service
[Unit]
Description=Mark new A/B slot as primary after successful boot
After=multi-user.target
Wants=multi-user.target
ConditionPathExists=/boot/syslinux.cfg

[Service]
Type=oneshot
ExecStart=/usr/sbin/boot-ok.sh

[Install]
WantedBy=multi-user.target
